# debian/mini-buildd-rep.cron.daily
#!/bin/sh -e

su - mini-buildd -c "/usr/lib/mini-buildd/mbd-update-rep"
